# systemd/snake-runner.service
[Unit]
Description=Snake Runner

[Service]
ExecStart=/usr/bin/snake-runner
Restart=always

[Install]
WantedBy=basic.target
